# dockers/configs/lava-test.service
[Unit]
Description=Lava test service

[Service]
ExecStart=/setup.sh

[Install]
WantedBy=multi-user.target
